# apps/terminal/deploy/systemd/terminal-stack.service
[Unit]
Description=Terminal Demo WebSocket Service Stack (Docker Compose)
Documentation=https://github.com/russellbomer/portfolio
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/terminal/deploy

# Environment file for secrets
# MANUAL: Create /var/www/terminal/.env with TOKEN_SECRET
EnvironmentFile=/var/www/terminal/.env

# Use docker compose to manage the stack
ExecStartPre=/usr/bin/docker compose pull --quiet
ExecStart=/usr/bin/docker compose up --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose restart

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terminal-stack

# Timeout settings
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
